# com.feishu.bot.supervisor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.feishu.bot.supervisor</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PROJECT_DIR__/venv/bin/supervisord</string>
        <string>-c</string>
        <string>__PROJECT_DIR__/supervisord.conf</string>
        <string>-n</string> 
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/logs/launchd.out</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/logs/launchd.err</string>
</dict>
</plist>
